# frontend/scripts/devops-dashboard-ui.service
[Unit]
Description=DevOps Dashboard PySide6 QML UI - labwc kiosk
After=systemd-user-sessions.service

[Service]
User=dragos
WorkingDirectory=/opt/devops-dashboard-ui

PAMName=login
TTYPath=/dev/tty1
TTYReset=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal

Environment=HOME=/home/dragos
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=XDG_CONFIG_HOME=/home/dragos/.config
Environment=QT_QPA_PLATFORM=wayland
Environment=QT_QUICK_CONTROLS_STYLE=Basic

ExecStart=/usr/bin/labwc -s "/opt/devops-dashboard-ui/.venv/bin/python3 /opt/devops-dashboard-ui/main.py"

Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target